# systemd/bot-sui.service
[Unit]
Description=Bot S-UI Telegram service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/bot-sui
Environment=PYTHONUNBUFFERED=1
Environment=BOT_SUI_BASE_DIR=/opt/bot-sui
Environment=BOT_SUI_CONFIG=/opt/bot-sui/config.json
ExecStart=/opt/bot-sui/venv/bin/python /opt/bot-sui/bot.py
Restart=always
RestartSec=5
User=bot-sui
Group=bot-sui
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/bot-sui

[Install]
WantedBy=multi-user.target
