# systemd/bot-sui-backup.service
[Unit]
Description=Bot S-UI backup service

[Service]
Type=oneshot
ExecStart=/root/bot-sui/backup.sh
User=root